# config/autostart/background.desktop
[Desktop Entry]
Exec=$HOME/.fehbg
Name=background
Type=Application
Version=1.0
