# config/autostart/Browser.desktop
[Desktop Entry]
Exec=/usr/bin/google-chrome-stable
Name=Browser
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
