# config/autostart/DisableMouseWhenOnKeyboard.desktop
[Desktop Entry]
Exec=syndaemon -i 2.5s -R -K -d -t
Name=DisableMouseWhenOnKeyboard
Type=Application
Version=1.0
