# config/autostart/opera.desktop
[Desktop Entry]
Exec=/usr/bin/opera
Name=opera
Type=Application
Version=1.0
